# ruantiblock-mod-lua/files/etc/uci-defaults/ruantiblock-mod-lua
uci set ruantiblock.config.bllist_module='/usr/libexec/ruantiblock/ruab_parser.lua'; uci commit ruantiblock
