# ruantiblock-mod-py/files/etc/uci-defaults/ruantiblock-mod-py
uci set ruantiblock.config.bllist_module='/usr/libexec/ruantiblock/ruab_parser.py'; uci commit ruantiblock
